# SolarProbe-Server/crontab
*/3 * * * *     /bin/bash -l -c '/var/lib/slang/library/SolarProbe-Server/render.sh'
